# Documents/Extra/ROOT_DIR/etc/profile.d/customvars.sh
export ALSA_CARD=PCH
export PATH="$PATH:/home/adigen/.local/bin"
export EDITOR="/usr/bin/nvim"
