# Documents/Extra/ROOT_DIR/etc/profile.d/customvars.sh
# Shifted this task to ~/.asoundrc
#export ALSA_CARD=PCH

export PATH="/home/adigen/.local/bin:$PATH"
export EDITOR="/usr/bin/nvim"
